# playbooks/roles/firewall/templates/vagrant_resolve.service.j2
[Unit]
Description=overwrite vagrant's attempts at changing dns

[Service]
ExecStart=/bin/sh /usr/sbin/vagrant_resolve.sh

[Install]
WantedBy=multi-user.target
