# playbooks/roles/firewall/templates/vagrant_resolve.service.j2
[Unit]
Description=overwrite vagrant's attempts at changing dns

[Service]
Type=simple
ExecStart=/bin/sh /usr/sbin/vagrant_resolve.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
